# example/akebi-keyless-server.socket
[Unit]
Description = Akebi Keyless Server Socket

[Socket]
ListenStream = 443
ListenDatagram = 53

[Install]
WantedBy = sockets.target
